# systemd/audiomuse-gpu-watch.timer
[Unit]
Description=Poll AudioMuse queue activity for GPU arbitration

[Timer]
OnBootSec=30s
OnUnitActiveSec=15s
Unit=audiomuse-gpu-watch.service

[Install]
WantedBy=timers.target
